# etc/cron.d/namecoin-zones
*/5 * * * * root /opt/namecoin-zones/update-bind.sh > /dev/null
